# ansible/roles/redis/templates/redis-server.service.j2
[Service]
# Redis systemd service override
# Managed by Ansible

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ redis_dir }}
ReadWritePaths=/var/log/redis

# Resource limits
LimitNOFILE=10032
